# etc/cron.daily/ntpclient
#!/bin/sh
bin=$(which ntpclient)
${bin} -s -h pool.ntp.org &>/dev/null
